# vpn-management-system/docker/update-agent/update-agent.service
[Unit]
Description=Update Agent for VPN Management System
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/vpn-management/update-agent
Environment="UPDATE_AGENT_TOKEN=changeme-update-token"
Environment="UPDATE_AGENT_PORT=8102"
Environment="INSTALL_DIR=/opt/vpn-management"
Environment="REPO_DIR=/opt/vpn-management/repo"
Environment="STATE_DIR=/var/lib/vpn-update"
Environment="UPDATE_SCRIPT=/opt/vpn-management/update-agent/update.sh"
# Bind on all interfaces so the backend container can reach it via host-gateway.
ExecStart=/opt/vpn-management/update-agent/venv/bin/gunicorn -w 2 -t 300 -b 0.0.0.0:8102 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
